# system/systemd/unlocator-monitor.service
[Unit]
Description=Unlocator VPN/DNS Reliability Monitor
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/unlocator-monitor.sh
