# system/systemd/unlocator-monitor.timer
[Unit]
Description=Unlocator Monitor - Every 5 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
