# monitoring/health-dashboard/health-dashboard.service
[Unit]
Description=SmartHome Health Dashboard
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/massey/health-dashboard
ExecStart=/usr/bin/python3 /home/massey/health-dashboard/app.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
